# OpenWrt/imagebuilder/files/etc/init.d/call_home_zerotier
#!/bin/sh /etc/rc.common

# Method to try and connect to a zerotier network

START=99

start() {

}